# timers/photos@.timer
[Unit]
Description=Run photos@%i every day

[Timer]
OnCalendar=daily
Persistent=true
Unit=photos@%i.service

[Install]
WantedBy=timers.target
